# ytblock-refresh.service
[Unit]
Description=ytblock IP refresh service
Documentation=https://github.com/user/ytblock
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ytblockctl reload
StandardOutput=journal
StandardError=journal